# ftmiddle.service.in
[Unit]
Description=Helium radio adjustment ftmiddle

[Service]
ExecStart=@@DESTROOT@@/run.sh
WorkingDirectory=@@DESTROOT@@
Restart=always

[Install]
WantedBy=multi-user.target
